# debian/postinst
#load the module
modprobe rtc-ds1307

systemctl enable rpirtc.service
systemctl start rpirtc.service
