# debian/postinst
systemctl enable rpirtc.service
systemctl start rpirtc.service
