# debian/rpirtc.cron.d
#
# Regular cron jobs for the rpirtc package
#
*/5 * * * * root /usr/bin/rpirtc_update
