# etc/systemd/system/rpirtc.service
[Unit]
Description=RPi RTC clock controller

[Service]
ExecStart=/usr/bin/rpirtc_start
Type=oneshot

[Install]
WantedBy=multi-user.target
